# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:greg]
command=/bin/bash -c "exec /build/greg.script start"
